# systemd/awake.target
[Unit]
Description=User-level Awake Target
Conflicts=sleep.target

[Install]
WantedBy=default.target
