# systemd/sleep.target
[Unit]
Description=User-level Sleep Target
Conflicts=awake.target
